# crontab
0 */12 * * * cd /app && python app.py >> /dev/null 2>&1